# docker/rootfs/etc/s6-overlay/s6-rc.d/opencode-im-bot/run
#!/bin/bash
# OpenCode IM Bot service

# Redirect stderr to stdout
exec 2>&1

echo "[opencode-im-bot] Starting IM Bot service..."

# Wait for OpenCode to be ready
echo "[opencode-im-bot] Waiting for OpenCode API..."
HOST="${OPENCODE_HOST:-127.0.0.1}"
PORT="${OPENCODE_PORT:-4096}"
MAX_RETRIES=30
COUNT=0

while ! curl -sf "http://${HOST}:${PORT}/health" >/dev/null; do
    sleep 2
    COUNT=$((COUNT+1))
    if [ "$COUNT" -ge "$MAX_RETRIES" ]; then
        echo "[opencode-im-bot] Timeout waiting for OpenCode API at http://${HOST}:${PORT}"
        exit 1
    fi
done

echo "[opencode-im-bot] OpenCode API ready"

cd /app

# Run the IM bot
exec s6-setuidgid cloudify \
    python -m opencode_on_im run --config /data/config.yaml
