# docker/rootfs/etc/s6-overlay/s6-rc.d/ttyd/finish
#!/bin/bash
# Finish script for ttyd service

EXIT_CODE=$1
echo "[ttyd] Service stopped with exit code ${EXIT_CODE}"

# Exit codes:
# 0   = Normal exit
# 143 = SIGTERM (graceful shutdown)
# 130 = SIGINT (Ctrl+C)

if [ "${EXIT_CODE}" -eq 0 ] || [ "${EXIT_CODE}" -eq 143 ] || [ "${EXIT_CODE}" -eq 130 ]; then
    exit 0
fi

# Abnormal exit - wait before s6 restarts
echo "[ttyd] Abnormal exit, restarting in 3 seconds..."
sleep 3
