# docker/rootfs/etc/s6-overlay/s6-rc.d/ttyd/run
#!/command/execlineb -P
with-contenv

# ttyd Web Terminal service
# Provides browser-based terminal access

fdmove -c 2 1

foreground { echo "[ttyd] Starting Web Terminal on port ${TTYD_PORT:-7681}..." }

cd /data

# Run ttyd with bash shell
# Options:
#   -p: Port
#   -W: Writable (allow input)
#   -t: Theme settings
#   -T: Set terminal type
exec ttyd
    -p ${TTYD_PORT:-7681}
    -W
    -t fontSize=14
    -t fontFamily="DejaVu Sans Mono, monospace"
    -t theme={"background":"#1e1e1e","foreground":"#d4d4d4"}
    -T xterm-256color
    /bin/bash -l
