# docker/rootfs/etc/s6-overlay/s6-rc.d/opencode/run
#!/bin/bash
# OpenCode AI coding assistant service

# Redirect stderr to stdout
exec 2>&1

echo "[opencode] Starting OpenCode service..."

# Set working directory to /data for projects
cd /data

# Check if OpenCode is available
if ! command -v opencode >/dev/null 2>&1; then
    echo "[opencode] Error: opencode command not found"
    exit 1
fi

# Run OpenCode with HTTP server enabled
# Port 4096 is the default OpenCode API port
exec s6-setuidgid cloudify \
    opencode serve --port "${OPENCODE_PORT:-4096}" --host 127.0.0.1
