# docker/rootfs/etc/s6-overlay/s6-rc.d/opencode/run
#!/command/execlineb -P
with-contenv

# OpenCode AI coding assistant service
# Requires: Node.js environment, OpenCode installed

fdmove -c 2 1

foreground { echo "[opencode] Starting OpenCode service..." }

# Set working directory to /data for projects
cd /data

# Check if OpenCode is available
backtick -n OPENCODE_BIN {
    pipeline { which opencode } cat
}
importas -u OPENCODE_BIN OPENCODE_BIN

if { test -n "${OPENCODE_BIN}" }

# Run OpenCode with HTTP server enabled
# Port 4096 is the default OpenCode API port
exec s6-setuidgid cloudify
opencode serve --port ${OPENCODE_PORT:-4096} --host 127.0.0.1
