# docker/rootfs/etc/s6-overlay/s6-rc.d/opencode/run
#!/bin/bash
# OpenCode AI coding assistant service

# Redirect stderr to stdout
exec 2>&1

echo "[opencode] Starting OpenCode service..."

# Set working directory to /data for projects
cd /data

# Check if OpenCode is available
if ! command -v opencode >/dev/null 2>&1; then
    echo "[opencode] 'opencode' command not found. Assuming external OpenCode instance."
    echo "[opencode] Internal service disabled."
    # Sleep forever to prevent restart loop and allow bot to connect to external instance
    exec tail -f /dev/null
fi

# Run OpenCode with HTTP server enabled
# Port 4096 is the default OpenCode API port
exec s6-setuidgid cloudify \
    opencode serve --port "${OPENCODE_PORT:-4096}" --host 127.0.0.1
